# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:mysqld]
command=/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true